# ansible/roles/aptly/templates/aptly-api.service.j2
[Unit]
Description=Aptly REST API service

[Service]
TimeoutStartSec=0
ExecStart=/bin/bash -c 'aptly api serve'

[Install]
WantedBy=multi-user.target
