# ansible/roles/hail/templates/hail.service.j2
[Unit]
Description=Hail Manager

[Service]
ExecStart={{ hail_service_start }}
User={{ hail_service_user }}
{% for key, value in hail_service_environment_variables.items() %}
Environment={{ key }}={{ value }}
{% endfor %}

[Install]
WantedBy=multi-user.target
